# etc/systemd/system/sentenial-threat-monitor.service
# /etc/systemd/system/sentenial-threat-monitor.service
[Unit]
Description=Sentenial Threat Monitor AI Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple

# Run as a non-root user/group
User=beastmode7777
Group=sentenial

# Ensure runtime dirs exist
ExecStartPre=/usr/bin/install -d -m 0755 /var/lib/sentenial
ExecStartPre=/usr/bin/install -d -m 0755 /opt/sentenial_x_ai

# Working directory where your code lives
WorkingDirectory=/opt/sentenial_x_ai

# Load defaults, then override with env file if present
EnvironmentFile=-/etc/sentenial-threat-monitor.env
Environment=DB_PATH=/var/lib/sentenial/agent.db
Environment=AGENT_ID=sentenial-x-ai-bot
Environment=NETWORK_PORT=8000
Environment=PEERS=http://localhost:5001,http://localhost:5002

# The actual start command
ExecStart=/usr/bin/python3 /opt/sentenial_x_ai/agent.py

# Graceful shutdown
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=5s

# Send stdout/stderr to journald/syslog
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentenial-threat-monitor

[Install]
WantedBy=multi-user.target
